# conversor.service
[Unit]
Description=Monitoreo de la carpeta input en Dropbox
After=network.target

[Service]
Type=simple
User=usuario
Group=usuario
ExecStart=/bin/bash -c 'source /home/usuario/env/bin/activate && exec /usr/local/bin/conversor.sh'
Restart=always
WorkingDirectory=/home/usuario
Environment="VIRTUAL_ENV=/home/usuario/env"
Environment="PATH=/home/usuario/env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
# Security enhancements
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=/home/usuario/temp /home/usuario/output /home/usuario/error_log /home/usuario/completed
NoNewPrivileges=yes
ProtectHome=false

[Install]
WantedBy=multi-user.target
